# packages/edge-gateway/xrfonstream-egw-memwatch.service
[Unit]
Description=SAX Edge Gateway memory watchdog (userspace MemoryMax fallback)

[Service]
Type=oneshot
ExecStart=/opt/sax/edge-gateway/memory-watchdog.sh
